# flash-unlock/data/etc/udev/rules.d/99-initrd-flash-unlock.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/uevent/filters/flash-unlock"
